# meta-agl-drm-lease/recipes-graphics/drm-lease-manager/drm-lease-manager/drm-lease-manager.service
[Unit]
Description=DRM Lease Manager

[Service]
ExecStart=drm-lease-manager
RuntimeDirectory=drm-lease-manager
